# cron-service/crontab
*/5 * * * * curl -sf -H "Authorization: Bearer ${CRON_SECRET}" ${SOFIA_URL}/api/flows/cron >> /proc/1/fd/1 2>&1
0 */1 * * * curl -sf -H "Authorization: Bearer ${CRON_SECRET}" ${SOFIA_URL}/api/cron/run-scheduled >> /proc/1/fd/1 2>&1
0 12 * * * curl -sf -H "Authorization: Bearer ${CRON_SECRET}" ${SOFIA_URL}/api/cron/email-drip >> /proc/1/fd/1 2>&1
0 11 * * 1 curl -sf -H "Authorization: Bearer ${CRON_SECRET}" ${SOFIA_URL}/api/cron/weekly-digest >> /proc/1/fd/1 2>&1
